# root/etc/s6-overlay/s6-rc.d/celery-beat/run
#!/command/with-contenv bash

export PATH=$S6_GLOBAL_PATH
s6-setuidgid nautobot nautobot-server celery beat -l $NAUTOBOT_LOG_LEVEL | while read line ; do
  echo "celery-beat: $line"
done
