# root/etc/s6-overlay/s6-rc.d/nautobot-worker/run
#!/command/with-contenv bash

export PATH=$S6_GLOBAL_PATH
s6-setuidgid nautobot nautobot-server rqworker | while read line ; do
  echo "nautobot-worker: $line"
done
