# root/etc/s6-overlay/s6-rc.d/nautobot-init/up
export PATH "/command:/usr/local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
export HOME "/opt/nautobot"
s6-setuidgid nautobot /docker-entrypoint.sh /install_plugins.sh
